# uxp/sonda.cron
# Run the ping_plugin (pingtest) every 10 minutes
*/2 * * * * cd /app && /usr/local/bin/python3 /app/main.py pingtest >> /dev/stdout 2>&1
*/4 * * * * cd /app && /usr/local/bin/python3 /app/main.py dnstest >> /dev/stdout 2>&1
